# files/systemd/user/session-init.service
[Service]
Type=oneshot
ExecStart=%h/.local/bin/session-init

[Install]
WantedBy=default.target
 
